# apps/control-plane-stack/assets/wg-peers.service
[Unit]
Description=Regenerate WireGuard peers from DB

[Service]
Type=oneshot
EnvironmentFile=/home/admin/control-plane/.env
ExecStart=/usr/local/bin/regen-wg-peers
